# asset-task-sync-backend/deploy/supervisor/laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/asset-task-sync-backend/artisan queue:work redis --sleep=3 --tries=3 --timeout=90
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/laravel-worker.log
stopwaitsecs=3600
priority=1000

; Notes:
; - Replace `/path/to/asset-task-sync-backend` with the absolute path on your server.
; - Ensure the `www-data` user exists and can run the php binary and write to the logfile.
; - After placing this file at `/etc/supervisor/conf.d/laravel-worker.conf` run:
;     sudo supervisorctl reread
;     sudo supervisorctl update
;     sudo supervisorctl start laravel-worker:*
